# scripts/launchd/com.codex-web-term.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex-web-term</string>

    <key>ProgramArguments</key>
    <array>
      <string>__NODE_PATH__</string>
      <string>__ROOT__/scripts/service.mjs</string>
      <string>resurrect</string>
      <string>--mode</string>
      <string>prod</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ROOT__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__ROOT__/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__ROOT__/logs/launchd-stderr.log</string>
  </dict>
</plist>
